# lib/football.conf
# Put this file in /etc/init/*.conf
# Usage - sudo start node
#		  sudo stop node.
# Easy peasey

description "Node service"
author "Nick Konecny - nkonecny.com" 

#
env NODE_PATH="/usr/local/bin/node"
env FULL_PATH="/home/pi/Github/BetTracker/"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown

# Automtically respawn
respawn                 # restart when job dies
respawn limit 99 5      # give up restart after 99 respawns in 5 seconds

script
	export HOME="/root"
	echo $$ > /var/run/$PROGRAM_NAME.pid
	cd $FULL_PATH
	exec sudo -u root $NODE_PATH $FULL_PATH/server.js >> /home/pi/logs/bettracker.log 2>&1
end script